# etc/profile.d/oj-welcome.sh
printf  "welcome $(cat /etc/issue | head -n1)\n"
